# scripts/systemd/mnemonic-train.path
[Unit]
Description=Watch for Mnemonic training requests

[Path]
PathExists=%h/.mnemonic/training_requests/pending.json
Unit=mnemonic-train.service

[Install]
WantedBy=default.target
